# miurd.service
[Unit]
Description=Miurd Discord Bot
After=network.target

[Service]
Type=simple
User=akaaku
WorkingDirectory=/path/to/miurd
ExecStart=/usr/bin/pnpm prod
Restart=always
RestartSec=10
Environment=NODE_ENV=production
# ログ設定
StandardOutput=append:/path/to/miurd/logs/systemd-out.log
StandardError=append:/path/to/miurd/logs/systemd-error.log

# セキュリティ設定
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/path/to/miurd

[Install]
WantedBy=multi-user.target